# rootfs/etc/services.d/daikin-smart-ac-br/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Preparing to start..."

# Read options
export DAIKINBR_CONFIG_RELOAD_THINGS=$(bashio::config 'reload_things')
export DAIKINBR_CONFIG_RELOAD_ADDRESSES=$(bashio::config 'reload_addresses')
export DAIKINBR_CONFIG_ACCOUNT_EMAIL=$(bashio::config 'daikin_email')
export DAIKINBR_CONFIG_ACCOUNT_PASSWORD=$(bashio::config 'daikin_password')


# Read MQTT configuration
if bashio::var.true "$(bashio::services 'mqtt' 'ssl')"; then
    export DAIKINBR_CONFIG_MQTT_SERVER="ssl://$(bashio::services 'mqtt' 'host'):$(bashio::services 'mqtt' 'port')"
else
    export DAIKINBR_CONFIG_MQTT_SERVER="tcp://$(bashio::services 'mqtt' 'host'):$(bashio::services 'mqtt' 'port')"
fi
export DAIKINBR_CONFIG_MQTT_USER="$(bashio::services 'mqtt' 'username')"
export DAIKINBR_CONFIG_MQTT_PASSWORD="$(bashio::services 'mqtt' 'password')"


# Starts server
bashio::log.info "Starting Daiking Smart AC Brasil..."
cd /app
./daikin-server
